# docker/raspberrypi/systemd/pg-backup.service
[Unit]
Description=BikeRental PostgreSQL backup (pg_dump)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# __USER__ and __PROJECT_DIR__ are filled in automatically by
# systemd/install-backup-timer.sh (it detects the user and the repo location).
# Running as your user (not root) keeps the dump files owned by you and uses your
# docker group membership. The backup DESTINATION is BACKUP_DIR, read from the
# EnvironmentFile (.env) at run time — change the destination there, not here.
User=__USER__
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=__PROJECT_DIR__/.env
ExecStart=__PROJECT_DIR__/backup/pg-backup.sh
