# docker/raspberrypi/systemd/pg-backup.service
[Unit]
Description=BikeRental PostgreSQL backup (pg_dump)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# __USER__ and the paths below are filled in automatically by
# systemd/install-backup-timer.sh. Running as your user (not root) keeps the
# dump files owned by you and uses your docker group membership.
User=__USER__
WorkingDirectory=/home/__USER__/bike-rental/docker/raspberrypi
EnvironmentFile=/home/__USER__/bike-rental/docker/raspberrypi/.env
ExecStart=/home/__USER__/bike-rental/docker/raspberrypi/backup/pg-backup.sh
